# .config/systemd/user/auto-commit-task-db.timer
[Unit]
Description=Commit task database twice an hour

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
