# .config/systemd/user/bugwarrior-pull.service
[Unit]
Description=bugwarrior-pull

[Service]
ExecStart=/usr/bin/bugwarrior-pull
Type=oneshot

[Install]
WantedBy=default.target
