# .config/systemd/user/clipboard_manager.service
[Unit]
Description=Clipboard manager

[Service]
ExecStart=/usr/bin/greenclip daemon
Restart=always

[Install]
WantedBy=default.target
