# .config/systemd/user/download-archiver.service
[Unit]
Description=Download Archiver

[Service]
ExecStart=%h/.bin/download-archiver
Type=oneshot

[Install]
WantedBy=default.target
