# .config/systemd/user/download-archiver.timer
[Unit]
Description=Archive Downloads monthly

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
