# .config/systemd/user/download-archiver.timer
[Unit]
Description=Archive Downloads weekly

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
