# .config/systemd/user/dunst.service
[Unit]
Description=Notification daemon

[Service]
ExecStart=/usr/bin/dunst
Restart=always

[Install]
WantedBy=default.target

