# .config/systemd/user/gwall.service
[Unit]
Description=Gwall

[Service]
ExecStart=%h/.bin/gwall

[Install]
WantedBy=default.target
