# .config/systemd/user/gwall.service
[Unit]
Description=Gwall

[Service]
ExecStart=/home/kalbasit/.bin/Gwall

[Install]
WantedBy=default.target
