# .config/systemd/user/hibernate-on-low-battery.timer
[Unit]
Description=Check for low battery every minute

[Timer]
OnBootSec=15min
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
